# scripts/systemd/orangecat-cat-eval.service
# Nightly Cat judgment regression gate.
# Runs scripts/eval-cat.mjs (deployed to /opt/orangecat/scripts/ by
# scripts/deploy-selfhost.sh) against production. On regression the script
# exits non-zero, logs to this unit's journal, and inserts an in-app
# notification for the founder (CAT_EVAL_NOTIFY=1).
#
# SSOT: this file lives in the repo at scripts/systemd/ and is installed to
# /etc/systemd/system/ by deploy-selfhost.sh — edit it here, not on the box.
[Unit]
Description=OrangeCat nightly Cat eval (entity-type judgment gate)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/opt/orangecat/app/.env
Environment=CAT_EVAL_NOTIFY=1
ExecStart=/usr/bin/node /opt/orangecat/scripts/eval-cat.mjs
