# scripts/systemd/orangecat-cat-eval.timer
# SSOT: repo scripts/systemd/ — installed by deploy-selfhost.sh.
[Unit]
Description=OrangeCat nightly Cat eval (04:30 UTC)

[Timer]
OnCalendar=*-*-* 04:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
